# files/rhythmbot.service
[Unit]
Description=Discord.js Bot Service
After=network.target

[Service]
# User running the bot (e.g., node, ubuntu, or root)
User=rhythmbot
WorkingDirectory=/home/rhythmbot/rhythmbot
ExecStart=/home/rhythmbot/.local/share/pnpm/pnpm run start
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rhythmbot

[Install]
WantedBy=multi-user.target